# launchd/com.assistant.todo-reminder.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
    Todo Reminder - launchd Service Template

    Periodically checks for overdue or high-priority to-dos and injects
    reminders into the Claude Code session.

    Installation:
    1. Copy to ~/Library/LaunchAgents/com.assistant.todo-reminder.plist
    2. Replace YOUR_USERNAME and YOUR_PROJECT with your values
    3. launchctl load ~/Library/LaunchAgents/com.assistant.todo-reminder.plist
    -->

    <key>Label</key>
    <string>com.assistant.todo-reminder</string>

    <key>ProgramArguments</key>
    <array>
        <!-- TODO: Update to your project directory -->
        <string>/Users/YOUR_USERNAME/YOUR_PROJECT/scripts/todo-reminder.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/YOUR_PROJECT</string>

    <key>StartInterval</key>
    <!-- Check every 30 minutes (1800 seconds) -->
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/YOUR_PROJECT/logs/todo-reminder.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/YOUR_PROJECT/logs/todo-reminder.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
