# launchd/com.assistant.email-reminder.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
    Email Reminder - launchd Service Template

    Periodically checks email inbox and injects unread count into the
    Claude Code session if there are new messages.

    Installation:
    1. Copy to ~/Library/LaunchAgents/com.assistant.email-reminder.plist
    2. Replace YOUR_USERNAME and YOUR_PROJECT with your values
    3. launchctl load ~/Library/LaunchAgents/com.assistant.email-reminder.plist
    -->

    <key>Label</key>
    <string>com.assistant.email-reminder</string>

    <key>ProgramArguments</key>
    <array>
        <!-- TODO: Update to your project directory -->
        <string>/Users/YOUR_USERNAME/YOUR_PROJECT/scripts/email-reminder.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/YOUR_PROJECT</string>

    <key>StartInterval</key>
    <!-- Check every 15 minutes (900 seconds) -->
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/YOUR_PROJECT/logs/email-reminder.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/YOUR_PROJECT/logs/email-reminder.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
